# systembus-notify.service
[Service]
ExecStart=%h/bin/systembus-notify

[Install]
WantedBy=default.target
